# containers/rootfs/etc/s6-overlay/s6-rc.d/crond/run
#!/command/execlineb -P
with-contenv
s6-notifyoncheck -d
/bin/crond -f -l 8 -L /dev/stdout -c /var/spool/cron/crontabs